# install/dbus/bus/dbus.socket
[Unit]
Description=D-Bus System Message Bus Socket

[Socket]
ListenStream=/home/chen/work/bluez/build/var/run/dbus/system_bus_socket
